# .docker/etc/s6/services/bluetooth/run
#!/bin/execlineb -WP
/usr/lib/bluetooth/bluetoothd -n
